# node-server.service
# Installation commands:
# sudo nano /etc/systemd/system/node-server.service
# sudo systemctl enable --now node-server
#
# Runs a local Node app as the unprivileged "pi" user (not root). This is set up 
# to run the bundled demo server in demo/. Point WorkingDirectory and ExecStart at
# your own project once you have one, and update User/Group if your account is
# named differently.

[Unit]
Description=Kiosk Node server
# Only orders after the network stack initialises; deliberately does NOT wait for
# network-online.target. The server binds localhost and serves the offline page
# itself, so it must start even with no internet at boot. Waiting for
# connectivity would delay startup by up to ~90s, during which Chromium hits 
# localhost too early and shows its own "site can't be reached" error instead 
# of our Connection lost page.
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/raspberry-pi-kiosk/demo
ExecStart=/usr/bin/node /home/pi/raspberry-pi-kiosk/demo/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nodeServer
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
